# crontabs
0 1 * * * bash /root/update-scripts/update-docker-compose-templates.sh
0 1 * * * bash /root/update-scripts/update-dockerdeploy.sh